# syntx-api.service
[Unit]
Description=SYNTX Production API v2 - Strom Orchestrator 🌊
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/syntx-workflow-api-get-prompts/api-core
ExecStart=/opt/syntx-workflow-api-get-prompts/venv/bin/python3 /opt/syntx-workflow-api-get-prompts/api-core/syntx_api_production_v2.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
